# config/stvmm-hypervisor/overlays/vanilla/usr/lib/systemd/system/stvmm-vsock@.service
[Unit]
Description=System Transparency Virtual Machine %I vsock server

[Service]
ExecStart=/usr/sbin/stvmmd receive

