# config/example/overlays/vanilla/usr/lib/systemd/system/st-example.service
[Unit]
Description=ST example service

[Service]
Type=exec
User=nobody
ExecStart=/usr/sbin/st-example-service

[Install]
WantedBy=multi-user.target
