# config/example/overlays/vanilla/usr/lib/systemd/system/st-init.service
[Unit]
Description=Initialise ST-booted system from /stdata
RequiresMountsFor=/stdata
Before=network.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/st-init

[Install]
WantedBy=basic.target
